# dots/crontab
* * * * * $HOME/custom/scripts/cron-rmdir-maj-home
* * * * * $HOME/custom/scripts/cron-partition-check
