# dots/crontab
*/5  * * * * $HOME/custom/scripts/cron/rmdir-maj-home
*/30 * * * * $HOME/custom/scripts/cron/partition-check
*/30 * * * * $HOME/custom/scripts/cron/check-update
